# etc/systemd/home-warden.service
[Unit]
Description=home-warden nginx (unprivileged)
Documentation=https://github.com/the-hcma/home-warden/blob/main/PLAN.md
Requires=home-warden.socket
After=network.target home-warden.socket

[Service]
Type=simple
User=@@OWNER@@
Group=@@OWNER@@
Environment=NGINX=3:4;
Environment=DEPLOYED_COMMIT=@@DEPLOYED_COMMIT@@
# -p sets prefix for relative pid/temp/log paths in the home nginx.conf.
ExecStart=/usr/sbin/nginx -p @@SCRATCH_DIR@@/ -c @@HOME_NGINX_CONF@@ -g 'daemon off;'
ExecReload=/usr/sbin/nginx -p @@SCRATCH_DIR@@/ -c @@HOME_NGINX_CONF@@ -s reload
Restart=on-failure
RestartSec=5
# Bind privileged ports only via inherited fds from home-warden.socket.
NoNewPrivileges=true
# Systemd / process logs (nginx access/error also under @@SCRATCH_DIR@@/logs/).
# Follow with: tail --follow=name --retry @@SCRATCH_DIR@@/home-warden.log
StandardOutput=append:@@SCRATCH_DIR@@/home-warden.log
StandardError=append:@@SCRATCH_DIR@@/home-warden.log

[Install]
WantedBy=multi-user.target
